# led-display.service
[Unit]
Description=LED Display Application with Auto Updates
After=network.target

[Service]
Type=simple
User=led-board
WorkingDirectory=/home/led-board/Desktop/led-board-project
ExecStart=/home/led-board/Desktop/led-board-project/venv/bin/python /home/led-board/Desktop/led-board-project/update_and_run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/led-board/Desktop/led-board-project

[Install]
WantedBy=multi-user.target 